# mediavault-mount-health.timer
[Unit]
Description=MediaVault NAS Mount Health Check Timer
Requires=mediavault-mount-health.service

[Timer]
# Run every 5 minutes
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
